# rospi-roscore.service.in
[Unit]
Description=ROS core server
After=network.target

[Service]
Type=simple
ExecStart=@ROS_DIR@/env.sh roscore
Restart=always
User=roscore
WorkingDirectory=~
StandardInput=null

[Install]
WantedBy=multi-user.target
